# tor/rootfs/etc/services.d/tor/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Tor
# Runs the Tor daemon
# ==============================================================================
bashio::log.info "Starting Tor..."

# Run the Tor daemon
exec tor
